# files/puppet.cron
# MANAGED BY PUPPET!
0,15,30,45 * * * * root /etc/puppet/apply.sh > /dev/null 2>&1
